# ansible/roles/ec2/files/etc/init/mount-ephemeral.conf
# mount-ephemeral
#
# mount swap and ephemeral disks on startup

description     "mount swap and ephemeral disks"

start on stopped rc RUNLEVEL=[2345]

task

exec /opt/ec2/mount-ephemeral.sh
